# .config/systemd/user/org-repo-update.service
[Unit]
Description=perform updates in the org repo

[Service]
ExecStart=%h/scripts/org-repo-update
